# _init/geodesic/rootfs/etc/profile.d/functions.sh
#!/bin/bash

echo "[sourced] etc/profile.d/functions.sh"

function hg() {
    history | grep "$@"
}
